# board/mathworks/x86/fs-overlay/etc/profile.d/prompt.sh
#!/bin/sh

# set the prompt
export PS1="\e[0;31mmwx86\> \e[m"
