# board/mathworks/zynqmp/boards/zcu111/fs-overlay/etc/profile.d/prompt.sh
#!/bin/sh

# set the prompt
export PS1="\[\e[0;31m\]zcu111\> \[\e[m\]"
