# board/mathworks/zynq/fs-overlay/etc/profile.d/prompt.sh
#!/bin/sh

# set the prompt
export PS1="\e[0;31mzynq\> \e[m"
